# rootfs/etc/crontabs/drupal
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
   *      *      *        *       *       drush --root /drupal/www cron-run
